# configuration/system/canup.service
[Unit]
Description=CAN ISOTP UP

[Service]
ExecStart=/opt/canup/canup.sh

[Install]
WantedBy=multi-user.target
