# configuration/system/canup.service
[Unit]
Description=CAN ISOTP UP

[Service]
ExecStart=/home/mfgtest/projects/towercontroller/canup.sh

[Install]
WantedBy=multi-user.target
